# config/supervisor/proxbet-telegram-bot.conf
[program:proxbet-telegram-bot]
command=php /var/www/html/backend/telegram_bot.php
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/proxbet-telegram-bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/proxbet-telegram-bot-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopwaitsecs=30
stopsignal=TERM
priority=999
